# files/etc/update-motd.d/00-header
#!/usr/bin/env bash

# Get distrib info (will set the $DISTRIB_DESCRIPTION variable)
[ -r /etc/lsb-release ] && . /etc/lsb-release

# Print the friggin' motd
echo "
                            ==(W{==========-      /===-
                              ||  (.--.)         /===-_---~~~~~~~~~------____
                              | \\_,|**|,__      |===-~___                _,-'
                 -==\\\\        \`\\ ' \`--'   ),    \`//~\\\\   ~~~~\`---.___.-~~
             ______-==|        /\`\\_. .__/\\ \\    | |  \\\\           _-~\`
       __--~~~  ,-/-==\\\\      (   | .  |~~~~|   | |   \`\\        ,'
    _-~       /'    |  \\\\     )__/==0==-\\<>/   / /      \\      /
  .'        /       |   \\\\      /~\\___/~~\\/  /' /        \\   /'
 /  ____  /         |    \\\`\\.__/-~~   \\  |_/'  /          \\/'
/-'~    ~~~~~---__  |     ~-/~         ( )   /'        _--~\`
                  \\_|      /        _) | ;  ),   __--~~
                    '~~--_/      _-~/- |/ \\   '-~ \\
                   {\\__--_/}    / \\\\_>-|)<__\\      \\
                   /'   (_/  _-~  | |__>--<__|      |
                  |   _/) )-~     | |__>--<__|      |
                  / /~ ,_/       / /__>---<__/      |
                 o-o _//        /-~_>---<__-~      /
                 (^(~          /~_>---<__-      _-~
                ,/|           /__>--<__/     _-~
             ,//('(          |__>--<__|     /                  .----_
            ( ( '))          |__>--<__|    |                 /' _---_~\\
         \`-)) )) (           |__>--<__|    |               /'  /     ~\\\`\\
        ,/,'//( (             \\__>--<__\\    \\            /'  //        ||
      ,( ( ((, ))              ~-__>--<_~-_  ~--____---~' _/'/        /'
    \`~/  )\` ) ,/|                 ~-_~>--<_/-__       __-~ _/
  ._-~//( )/ )) \`                    ~~-'_/_/ /~~~~~~~__--~
   ;'( ')/ ,)(                              ~~~~~~~~~~
  ' ') '( (/
    '   '  \`
"
printf "Hey, welcome to %s (%s %s %s).\n\n" "${DISTRIB_DESCRIPTION}" "$(uname -o)" "$(uname -r)" "$(uname -m)"
printf "Just be nice, don't mess with the dragon-mounted gatekeeper above, and everything will be fine.\n\n"
